# weatherfeed.service
[Unit]
Description=Weatherfeed — hourly weather check and Telegram notifications
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/weatherfeed/venv/bin/python3 /opt/weatherfeed/runner.py
